# build/configs/standalone/supervisor/ds-docservice.conf
[program:docservice]
command=/var/www/euro-office/documentserver/server/DocService/docservice
directory=/var/www/euro-office/documentserver/server/DocService
user=ds
environment=NODE_ENV=production-linux,NODE_CONFIG_DIR=/etc/euro-office/documentserver,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice,LD_LIBRARY_PATH=/var/www/euro-office/documentserver/server/FileConverter/bin:/var/www/euro-office/documentserver/server/FileConverter/lib
stdout_logfile=/var/log/euro-office/documentserver/docservice/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/euro-office/documentserver/docservice/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
