# build/configs/standalone/supervisor/ds-docservice.conf
[program:docservice]
command=%(ENV_EO_ROOT)s/server/DocService/docservice
directory=%(ENV_EO_ROOT)s/server/DocService
user=ds
environment=NODE_ENV=production-linux,NODE_CONFIG_DIR=%(ENV_EO_CONF)s,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice,LD_LIBRARY_PATH=%(ENV_EO_ROOT)s/server/FileConverter/bin:%(ENV_EO_ROOT)s/server/FileConverter/lib
stdout_logfile=%(ENV_EO_LOG)s/docservice/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_EO_LOG)s/docservice/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
