# build/configs/standalone/supervisor/ds-example.conf
[program:ds-example]
command=%(ENV_EO_ROOT)s-example/example
directory=%(ENV_EO_ROOT)s-example/
user=ds
environment=NODE_ENV=production-linux,NODE_CONFIG_DIR=%(ENV_EO_CONF)s-example,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice
stdout_logfile=%(ENV_EO_LOG)s/ds-example_out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_EO_LOG)s/ds-example_err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=false
