# build/configs/standalone/supervisor/ds-metrics.conf
[program:metrics]
command=%(ENV_EO_ROOT)s/server/Metrics/metrics ./config/config.js
directory=%(ENV_EO_ROOT)s/server/Metrics
user=ds
environment=NODE_DISABLE_COLORS=1,LD_LIBRARY_PATH=%(ENV_EO_ROOT)s/server/FileConverter/bin:%(ENV_EO_ROOT)s/server/FileConverter/lib
stdout_logfile=%(ENV_EO_LOG)s/metrics/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_EO_LOG)s/metrics/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
